# dist/debian/postinst
#!/bin/sh
# Run on install AND on every upgrade. Idempotent — checks for prior state
# before mutating. Mirrors the dmart.spec %post block (RPM-side) so behavior
# stays consistent across package formats.
set -e

case "$1" in
    configure)
        # System user/group for the daemon. Created on first install only;
        # idempotent on upgrades because adduser --system is a no-op when
        # the account already exists.
        if ! getent group dmart >/dev/null; then
            addgroup --system dmart
        fi
        if ! getent passwd dmart >/dev/null; then
            adduser --system --no-create-home \
                    --home /var/lib/dmart --shell /usr/sbin/nologin \
                    --ingroup dmart --gecos "dmart data platform" dmart
        fi

        # /etc/dmart/config.env ships in the .deb as a conffile (see
        # dist/build-deb.sh — initial copy of config.env.sample). dpkg
        # installs it on first install with root:root 0640 from the
        # package tree, then prompts the operator on upgrade if it
        # diverges. We only need to chown to root:dmart so the daemon
        # (running as dmart) can read it.
        if [ -f /etc/dmart/config.env ]; then
            chown root:dmart /etc/dmart/config.env 2>/dev/null || true
            chmod 0640 /etc/dmart/config.env 2>/dev/null || true
        fi
        # Friendly first-install hint (skipped on upgrades — $2 is the
        # previously-installed version on upgrade, empty on fresh install).
        if [ -z "${2:-}" ]; then
            cat <<EOF
Installed default config: /etc/dmart/config.env
Edit it with your database credentials and JWT secret, then:
  systemctl enable --now dmart
EOF
        fi

        # Re-create runtime directories if missing, then ensure they're
        # owned by the dmart user. Belt + suspenders: the .deb's tar
        # entries SHOULD ship these, and the systemd unit's
        # StateDirectory= will recreate them at service start — but an
        # unclean upgrade or manual cleanup between install and first
        # start can still leave the dmart user staring at ENOENT. install
        # -d is idempotent (no-op when the directory exists with the
        # right perms) so this also acts as a self-heal on `dpkg
        # --reconfigure dmart`.
        install -d -m 0755 -o dmart -g dmart /var/lib/dmart
        install -d -m 0755 -o dmart -g dmart /var/lib/dmart/spaces
        install -d -m 0755 -o dmart -g dmart /var/lib/dmart/custom_plugins
        install -d -m 0755 -o dmart -g dmart /var/lib/dmart/logs

        # Reload systemd if it's there (containers / minimal images
        # might not have it; don't fail the install on its absence).
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
