# admin_scripts/docker/dmart.sysv
#!/sbin/openrc-run
supervisor=supervise-daemon
name="DMART Backend service"
description="Structured CMS/IMS — C# port"

command="/usr/bin/dmart"
command_args="serve --cxb-config /root/.dmart/config.json"
command_background="true"

# config.env lookup
export BACKEND_ENV="/root/.dmart/config.env"

output_log="/home/logs/dmart-out"
error_log="/home/logs/dmart-err"

depend() {
  use net localmount logger
  want postgresql
  after postgresql
}
